# aptrepo/seed-aptsigner/fs/etc/cron.d/aptsigner-repo-sync
*/2 * * * * root  rsync -qauze "ssh -i <key>" --delete-after /vol/aptrepo-local/repositories/raspbian /vol/aptrepo-local/repositories/ubuntu -avz aptsigner@<repo>:/
