# aptrepo/seed-aptsigner/fs/etc/cron.d/aptsigner-poller
# Keep scripts updated
1 * * * * root curl https://raw.githubusercontent.com/open-horizon/horizon-deb-packager/master/tools/apt-poller -o /usr/local/bin/apt-poller
1 * * * * root curl https://raw.githubusercontent.com/open-horizon/horizon-deb-packager/master/tools/promote-debs -o /usr/local/bin/promote-debs
# Continually check for updates
*/5 * * * * root sleep 30; /usr/local/bin/apt-poller -a http://pkg.bluehorizon.network/linux/ubuntu/dists/xenial-testing/main/binary-amd64/Packages -s <KEY> -b master
1-59/5 * * * * root sleep 30; /usr/local/bin/apt-poller -a http://pkg.bluehorizon.network/linux/ubuntu/dists/xenial-la-testing/main/binary-amd64/Packages -s <KEY> -b la
2-59/5 * * * * root sleep 30; /usr/local/bin/apt-poller -a http://pkg.bluehorizon.network/linux/ubuntu/dists/xenial-v3.2.1-testing/main/binary-amd64/Packages -s <KEY> -b v3.2.1
