# pkgsrc/deb/shared/debian/bluehorizon.postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable horizon.service
systemctl start horizon.service

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#


exit 0
