# pkgsrc/deb/shared/debian/bluehorizon.postinst
#!/usr/bin/env bash

systemctl daemon-reload
systemctl enable horizon.service
systemctl start horizon.service

exit 0
