# pkgsrc/deb/shared/debian/horizon-wiotp.postinst
#!/bin/bash

ln -s /usr/wiotp-edge/bin/wiotp_create_certificate.sh /usr/bin/wiotp_create_certificate
ln -s /usr/wiotp-edge/bin/wiotp_agent_setup.sh /usr/bin/wiotp_agent_setup
ln -s /usr/wiotp-edge/bin/wiotp_gwrlac_grp_add_device.sh /usr/bin/wiotp_gwrlac_grp_add_device

# We can now start the horizon service b/c horizon-wiotp contains all of the
# necessary config (like bluehorizon does)
systemctl daemon-reload
systemctl enable horizon.service
systemctl start horizon.service

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
